# deploy/tail-link.service
[Unit]
Description=Tail.Link FastAPI Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tail-link
Group=tail-link
WorkingDirectory=/opt/tail-link/backend
EnvironmentFile=/etc/tail-link.env
# SQLite + 内存握手状态要求单 worker；Nginx 负责公网入口。
ExecStart=/opt/tail-link/backend/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000 --workers 1 --proxy-headers --forwarded-allow-ips=127.0.0.1
Restart=always
RestartSec=3
TimeoutStopSec=20
Environment=PYTHONUNBUFFERED=1
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ReadWritePaths=/var/lib/tail-link

[Install]
WantedBy=multi-user.target
